# dist-linux/rift.service
[Unit]
Description=Rift (Bedrock proxy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# ── Update path to match your deploy location ──
WorkingDirectory=/opt/rift
ExecStart=/opt/rift/rift-musl /opt/rift/config.toml
# Automatically restart on unexpected exit
Restart=on-failure
RestartSec=3
# Logs are sent to journald → journalctl -u rift -f
StandardOutput=journal
StandardError=journal
Environment=RUST_LOG=rift=info
# Raise the file-descriptor limit for high UDP connection counts
LimitNOFILE=65536

# Note: stdin is not available in systemd mode, so console commands (transfer/kick/stop) cannot be used.
#   → Manage via the web dashboard (:8080) and systemctl. Use screen/tmux if you need console commands.

[Install]
WantedBy=multi-user.target
